# opendomain.service
[Unit]
Description=OpenDomain - Free Domain Registration Platform
After=network.target postgresql.service redis.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/opendomain
EnvironmentFile=/var/www/opendomain/.env
ExecStart=/var/www/opendomain/opendomain
Restart=always
RestartSec=5

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/opendomain/logs

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=opendomain

[Install]
WantedBy=multi-user.target
